# init/openrc/nic
#!/sbin/openrc-run

description="nic - network interface configurator"
command="/usr/local/sbin/nic"
command_args="start --daemon"
pidfile="/run/nic/nic.pid"

depend() {
    need localmount
    before net dns
    after udev
    provide net
}

start() {
    ebegin "Starting nic"
    ${command} ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping nic"
    ${command} stop
    eend $?
}

reload() {
    ebegin "Reloading nic"
    ${command} reload --force
    eend $?
}
